# data/boot.timer
[Unit]
Description=Take snapper snapshot of root on boot

[Timer]
OnBootSec=1

[Install]
WantedBy=timers.target
